# images/scripts/test.thing-workarounds/tt-vsock-ssh-service
#!/sbin/openrc-run

command="/usr/local/sbin/tt-vsock-sshd"
command_background=true
pidfile=/run/tt-vsock-sshd.pid
